# systemd/.config/systemd/user/backup.service
[Unit]
Description=Run restic backup script

[Service]
Type=simple
ExecStart=%h/scripts/restic-backup.sh
